# project-atlantis/scoring-service/gunicorn.service
[Service]
WorkingDirectory=/app
ExecStart=/usr/local/bin/gunicorn -c /app/gunicorn.conf.py