# AI-tutorbot/contrib/dockerizing/scoring-service/container_settings/gunicorn.service
[Service]
WorkingDirectory=/app
ExecStart=/usr/local/bin/gunicorn --bind unix:/run/gunicorn.sock scoring_service.wsgi:application